# tl.service
[Unit]
Description=Tl

[Service]
ExecStart=$HOME/bin/tl start
Restart=always
